# dist/pihub/install
#!/bin/sh -xe

sudo apt-get update -y
sudo apt-get upgrade -y
sudo apt-get install python3-pip -y
sudo pip3 install adafruit-circuitpython-lis3dh
sudo pip3 install adafruit-circuitpython-ads1x15

install_location=${INSTALL_TO:-$PWD}

sudo tee /etc/systemd/system/pihub.service << EOF
[Unit]
Description=Pihub
After=network.target

[Service]
ExecStart=${install_location}/pihub
WorkingDirectory=${install_location}
StandardOutput=inherit
StandardError=inherit
Restart=always
User=${USER}

[Install]
WantedBy=multi-user.target
EOF

if sudo systemctl is-enabled pihub; then
	sudo systemctl daemon-reload
	sudo systemctl restart pihub
else
	sudo systemctl enable pihub
	sudo systemctl start pihub
fi

# check whether it's up
if sudo systemctl status --no-pager pihub; then
	echo << EOF
The pihub service has been successfully installed. It should now
be listening on 0.0.0.0:3141.
EOF
else
	echo << EOF
Something went wrong and pihub does not appear to be properly
installed. Check the logs by running

sudo journalctl -u pihub
EOF
fi
